# recipes-apps/king/files/king.service
[Unit]
Description=A simple application

[Service]
Type=forking
ExecStart=/usr/bin/king.py

[Install]
WantedBy=multi-user.target
